# elasticsearch.service.conf
[program:elasticsearch]
command=/opt/elasticsearch/bin/elasticsearch -Des.config=/storage/config/elasticsearch.yml
