# services/elasticsearch.service.conf
[program:elasticsearch]
user=root
command=/opt/elasticsearch/bin/elasticsearch -Des.config=/storage/config/elasticsearch.yml
stdout_events_enabled = true
stderr_events_enabled = true
