# cookbooks/elasticsearch/files/default/cron.service.conf
[program:crond]
user=root
command=/usr/sbin/crond -n -x ext
numprocs=1
autostart=true
autorestart=true
startsecs=10
startretries=3
exitcodes=0,2
stopsignal=TERM
stopwaitsecs=10
